# systemd/client/btrfs-sync-client.service
[Unit]
Description=btrfs sync client

[Service]
ExecStart=/opt/btrfs-sync-daemon/client.py
