# systemd/client/btrfs-sync-client.timer
[Unit]
Description=Run btrfs sync client at 02:00 nightly

[Timer]
OnCalendar=02:00:00

[Install]
WantedBy=timers.target
